# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-firstboot.service.d/10-snosi-never-prompt.conf
# systemd-firstboot is a STATIC unit (no [Install] section), so it cannot be
# preset-disabled. With real first-boot semantics (machine-id ships as
# "uninitialized") it would run on every install's first boot and can prompt
# on the console for anything unset. Locale/hostname/user setup is owned by
# the installer and first-setup — keep this unit permanently inert.
# (Adding a condition ANDs it with the stock ones; do NOT write an empty
# Condition*= reset here.)
[Unit]
ConditionPathExists=/nonexistent
