# shared/snow/tree/usr/lib/systemd/system/snow-linux-live-setup.service
# SPDX-License-Identifier: GPL-2.0-only OR GPL-3.0-only OR LicenseRef-KDE-Accepted-GPL
# SPDX-FileCopyrightText: 2023 Harald Sitter <sitter@kde.org>
# Creates the default user and performs initial setup for SNOW Linux live sessions.
[Unit]
Description=SNOW Linux Live Session Setup
Requires=systemd-logind.service
Before=display-manager.service
After=systemd-logind.service systemd-homed.service
ConditionPathExists=/usr/bin/_snow-linux-live-setup
# Positive live-media detection (docs/plans/2026-07-14-bootc-native-ab-
# coexistence-plan.md, "Snow Native A/B": this unit's old marker-file-only
# gate is a NEGATIVE heuristic -- "absence of the done marker" -- which is
# indistinguishable from a freshly installed native A/B system's first boot
# (fresh /var, no marker). That false-positive would create a passwordless
# sudo "snow" user on every real install. Mirror the SAME positive
# kernel-cmdline signal already used by docker.socket.d/override.conf,
# incus.socket.d/override.conf, and brew-setup.service (which skip
# themselves with `!snow-linux.live=1`); this unit is the inverse case (it
# must run ONLY on live media), so it requires the cmdline arg directly.
ConditionKernelCommandLine=snow-linux.live=1
# Run-once via marker file instead of `systemctl disable`: the disable
# removed this unit's .wants/.requires symlinks from /etc at runtime, and
# any path deleted from the live /etc makes bootc's /etc merge fail at
# update finalize (bootc <= 1.16.3 follows the corresponding symlink in the
# new deployment's /etc and aborts: "a path led outside of the filesystem").
# Live sessions are unaffected: their state is ephemeral, so the marker
# resets every live boot exactly like the old runtime disable did.
ConditionPathExists=!/var/lib/snow-linux-live-setup.done

[Service]
ExecStart=/usr/bin/_snow-linux-live-setup
ExecStartPost=/usr/bin/touch /var/lib/snow-linux-live-setup.done

[Install]
WantedBy=multi-user.target
RequiredBy=multi-user.target
RequiredBy=display-manager.service
