# shared/outformat/ab-root/tree/usr/lib/systemd/system/snosi-ask-password-serial.service
# Dispatch password requests to the serial console on desktop native A/B
# images. Their cmdline ends with console=tty0 (so plymouth's local console
# is the VT and the graphical splash works -- see the plymouth note in
# shared/composition/snow/mkosi.conf), which moves /dev/console to tty0 and
# takes systemd-ask-password-console's prompt OFF the serial port. This
# agent restores the serial prompt (first-boot LUKS recovery passphrase
# before TPM enrollment -- also what the QEMU harnesses' console pump
# types into). Deliberately NO plymouth condition: with
# plymouth.ignore-serial-consoles, plymouth never manages ttyS0, so this
# agent must serve serial even while plymouth is running. Conditions keep
# it off cayo (no console=tty0 there) and off hardware without a serial
# device. Shipped into the native initrd via install_items in
# shared/outformat/ab-root/tree/usr/lib/dracut/dracut.conf.d/30-bootc-standard.conf.
[Unit]
Description=Dispatch Password Requests to Serial Console (snosi)
ConditionKernelCommandLine=console=tty0
ConditionKernelCommandLine=console=ttyS0
ConditionPathExists=/dev/ttyS0
DefaultDependencies=no
Conflicts=emergency.service
Before=emergency.service
Conflicts=shutdown.target initrd-switch-root.target
Before=shutdown.target initrd-switch-root.target

[Service]
ExecStart=systemd-tty-ask-password-agent --watch --console=/dev/ttyS0
SystemCallArchitectures=native
