# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/sshd-keygen.service.d/10-run-when-keys-missing.conf
# mkosi ships /etc/machine-id as an empty file, which systemd treats as
# "initialize a machine ID, but this is not first boot" (first-boot semantics
# require the file to be missing or to contain "uninitialized"). The stock
# unit's ConditionFirstBoot=yes therefore never fires on installed systems,
# no host keys are ever generated, and sshd crash-loops.
#
# Run whenever host keys are missing instead. The empty ConditionFirstBoot=
# resets ALL conditions on the unit, so the stock unit's other conditions are
# restated below.
[Unit]
ConditionFirstBoot=
ConditionPathIsReadWrite=/etc/ssh
ConditionPathIsSymbolicLink=!/etc/ssh
ConditionPathExists=!/etc/ssh/ssh_host_ed25519_key
