# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/bootc-update-stage.timer
# SPDX-License-Identifier: LGPL-2.1-or-later
[Unit]
Description=Stage bootc update timer

[Timer]
# Mirrors nbc-update-download.timer's cadence during the transition.
OnCalendar=hourly
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
