# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/preset-global.service
# SPDX-License-Identifier: LGPL-2.1-or-later
# Apply user-scope unit presets on first boot, mirroring what PID 1 does for
# system units. The image ships NO enablement symlinks in /etc (see the
# outformat finalize script); this recreates the user-scope ones from
# /usr/lib/systemd/user-preset policy as runtime-created /etc state.
# Adapted from ParticleOS. Safe for the bootc /etc merge: it runs only on
# first boot, when /etc/systemd/user contains no enablement symlinks, so
# preset-all can only CREATE paths, never delete shipped ones.
[Unit]
Description=Apply user unit presets on first boot
ConditionFirstBoot=yes
ConditionPathIsReadWrite=/etc

DefaultDependencies=no

Before=basic.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
# etc-guard-allow: first-boot only; image /etc ships no enablement symlinks, so this only creates paths
ExecStart=systemctl preset-all --global

[Install]
WantedBy=basic.target
