# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/preset-migration.service
# SPDX-License-Identifier: LGPL-2.1-or-later
# One-shot migration to preset-owned enablement for EXISTING installs.
#
# Images ship no unit enablement symlinks in /etc (presets recreate them at
# first boot — see preset-global.service and the outformat finalize script).
# But an install created from a pre-2026-07 image already had its first
# boot, and bootc's 3-way /etc merge drops the old image's enablement
# symlinks on the first update across the model boundary (unchanged in the
# live /etc + absent from the new image /etc = removed). Without this unit
# the system would come up with every preset-enabled service inactive.
#
# Activation is a STATIC wants symlink in /usr (sysinit.target.wants), not
# an /etc enablement — after the boundary update there are no /etc
# enablement symlinks left to activate anything, and /usr wants dirs are
# always scanned. The conditions keep it inert everywhere else:
#   - fresh installs:  first boot handles presets (ConditionFirstBoot=no),
#     and preset-global.service writes the marker on first boot
#   - migrated system: marker exists after the one run
# Do NOT gate on multi-user.target.wants being empty: runtime-created or
# locally-modified enablement symlinks legitimately survive the /etc merge
# (observed: ipp-usb.service), so the dir is not reliably empty at the
# boundary. Deleting /var/lib/preset-enablement.done re-runs the presets
# and costs one reboot — don't.
#
# After applying presets it reboots once: the boot transaction was computed
# before the enablement symlinks existed, so the created state only takes
# effect on the next boot. Marker is written first — no reboot loop.
[Unit]
Description=Re-apply unit presets after migration to preset-owned enablement
DefaultDependencies=no
ConditionFirstBoot=no
ConditionPathExists=!/var/lib/preset-enablement.done
ConditionPathIsReadWrite=/etc
# DefaultDependencies=no opts out of the implicit ordering after local-fs;
# the marker below lives in /var, so wait for it explicitly or the touch
# races the /var mount and the unit fails before its reboot step.
RequiresMountsFor=/var/lib
After=local-fs.target
Before=sysinit.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
# etc-guard-allow: one-shot at the enablement-model boundary; /etc holds no enablement symlinks at this point, so preset-all only creates paths
ExecStart=systemctl preset-all
# etc-guard-allow: same justification, user scope
ExecStart=systemctl preset-all --global
# Marker content = durable evidence this path ran: the reboot below happens
# before systemd-journal-flush, so this boot's journal never reaches disk.
ExecStart=/bin/sh -c 'echo migration > /var/lib/preset-enablement.done'
ExecStart=systemctl --no-block reboot
