# mkosi.images/base/mkosi.extra/usr/lib/systemd/user/bootc-update-notify.service
# SPDX-License-Identifier: LGPL-2.1-or-later
# Per-user desktop notification for a staged OS update. Ack-gated per staged
# digest: notifies once per staged update, not per login or path trigger.
[Unit]
Description=Notify about a staged OS update pending reboot
After=graphical-session.target
ConditionPathExists=/run/snosi/update-staged
# See bootc-update-notify.path: one staging can trigger this several times in a
# burst. Repeat runs no-op via the per-digest ack but still count as starts, so
# disable the default start-rate limit to avoid a spurious start-limit failure.
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/libexec/bootc-update-notify

# No [Install]: triggered solely by bootc-update-notify.path (whose
# PathExists= also fires at session start when the semaphore already exists).
