# shared/snow/tree/usr/lib/systemd/system/plymouth-start.service.d/10-wait-drm.conf
# Wait (bounded) for the first DRM device before starting plymouthd.
# plymouthd enumerates DRM devices ONCE at startup; if the GPU driver has
# not created /dev/dri/card0 yet (probabilistic race, observed live in QEMU
# 2026-07-17: the same build flips between the graphical two-step splash
# and the text-splash fallback purely on boot timing), it permanently falls
# back to the text splash -- a later udev add event does not reliably
# upgrade an active splash. `udevadm wait` blocks until the devnode exists
# AND udev has initialized it (`udevadm settle` is NOT sufficient -- it
# returns early when the probe event has not been queued yet; same lesson
# as snosi-setup's abandoned DRM device-pull). The `-` prefix keeps a
# GPU-less or slow system booting: on timeout plymouth just starts in its
# old racy/details mode. Ships in the snow tree (desktop payload, both
# transports); harmless on bootc where initrd plymouthd is already running
# by the time this unit is considered.
[Service]
ExecStartPre=-/usr/bin/udevadm wait /dev/dri/card0 --timeout=10
